# srcs/mysql/supervisor.conf
[supervisord]
nodaemon = true

[program:mysqld]
command = sh -c "mysqld_safe && kill -s SIGTERM $(cat supervisord.pid)"