# srcs/grafana/supervisor.conf
[supervisord]
nodaemon = true

[program:grafana]
command = sh -c "./grafana_server && kill -s SIGTERM $(cat supervisord.pid)"